# jobs/scheduler/launchd/com.nandish.jobs.ats.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
        "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nandish.jobs.ats</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/nandish.chokshi/Downloads/portfolio/jobs/scheduler/run_ats.sh</string>
    </array>

    <!-- Run every hour, on the hour -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Don't pile up if a run overlaps -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/nandish.chokshi/Downloads/portfolio/jobs/data/logs/ats.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nandish.chokshi/Downloads/portfolio/jobs/data/logs/ats.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
